# debian/sensorfw.conf
description "Sensorfw Services"

start on started dbus

respawn

exec /usr/sbin/sensorfwd --device-info --log-level=warning --no-magnetometer-bg-calibration
